# deploy/systemd/vol-shield-streamlit.service
[Unit]
Description=Vol shield Streamlit dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/quant/vol-shield
EnvironmentFile=/root/quant/vol-shield/.env
ExecStart=/root/quant/vol-shield/.venv/bin/streamlit run app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
